# airootfs/etc/systemd/system/os-installer.service
[Unit]
Description=OS Installer

[Service]
ExecStart=/usr/bin/os-installer
User=root

[Install]
WantedBy=default.target
